# rootfs/etc/periodic/daily/update-ttrss
#!/usr/bin/execlineb -P
foreground {
cd /usr/html
exec -c git pull
}
s6-svc -h /var/run/s6/services/ttrss
